# etc/crontab
0 * * * * /usr/local/bin/python /app/jcr/manage.py scrape_races >> /app/var/log/cron.log 2>&1
